# systemd/claver-agent.service
[Unit]
Description=Claver VPS Agent
Documentation=https://github.com/rockclaver/claver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claver
Group=claver
ExecStart=/usr/local/bin/claver-agent --addr 127.0.0.1:7676
Restart=on-failure
RestartSec=2s

# Hardening.
# NoNewPrivileges MUST stay off: the agent runs as the unprivileged `claver`
# user and escalates privileged operations (firewall ufw/firewall-cmd from
# issue #44, and host reboot) through `sudo -n`, scoped by sudoers.d fragments
# to specific binaries/verbs. NoNewPrivileges blocks sudo's setuid, which
# breaks every one of those paths. Crucially, several systemd sandboxing
# directives (PrivateDevices=, ProtectKernelTunables=, ProtectKernelModules=,
# ProtectControlGroups=, RestrictAddressFamilies=) *implicitly* force
# NoNewPrivileges=yes (see systemd.exec(5)), so they are deliberately omitted
# here — setting them silently re-broke sudo despite the explicit intent.
# All privileged mutations still pass through the agent's biometric +
# single-use confirmation-token (+ firewall anti-lockout) pipeline before sudo.
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# The agent's *listener* is locked to loopback by an in-process self-check
# (assertLoopback). It still needs outbound network to clone repos, talk to
# GitHub CLI auth, previews, etc., so we leave egress open here.
ReadWritePaths=/var/lib/claver /etc/caddy/claver

[Install]
WantedBy=multi-user.target
