# systemd/claver-agent.service
[Unit]
Description=Claver VPS Agent
Documentation=https://github.com/rockclaver/claver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claver
Group=claver
ExecStart=/usr/local/bin/claver-agent --addr 127.0.0.1:7676
Restart=on-failure
RestartSec=2s

# Hardening.
# NoNewPrivileges is intentionally NOT set: firewall management
# (ufw / firewall-cmd, issue #44) requires root and the agent escalates
# via a sudoers.d fragment scoped to just those two binaries
# (/etc/sudoers.d/claver-agent-firewall). NoNewPrivileges blocks setuid
# which would break that path. All firewall mutations still pass through
# the agent's biometric + confirmation-token + anti-lockout pipeline
# before sudo is invoked.
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
# The agent's *listener* is locked to loopback by an in-process self-check
# (assertLoopback). It still needs outbound network to clone repos, talk to
# GitHub CLI auth, previews, etc., so we leave egress open here.
ReadWritePaths=/var/lib/claver /etc/caddy/claver

[Install]
WantedBy=multi-user.target
